# deploy/vpn-backup.timer
[Unit]
Description=Run VPN Manager backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
